# cattexts.conf
[program:cattexts]
command=/home/ubuntu/cat-facts-twilio/cattexts.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/cattexts.err.log
stdout_logfile=/var/log/cattexts.out.log
